# tools/api/api.service
[Unit]
Description=Vengi API
After=network.target

[Service]
Type=forking
ExecStart=vengi.py

[Install]
WantedBy=multi-user.target
