# contrib/installer/linux/server.service.in
[Unit]
Description=@DESCRIPTION@
After=network.target

[Service]
Type=forking
ExecStart=@CMAKE_PROJECT_NAME@-@_EXE_TARGET@

[Install]
WantedBy=multi-user.target
